# supervisord.conf
[supervisord]

nodaemon=true

[program:mysql]

command=/start.sh mysqld

[program:rails]

command=rails server -p 80 -d
